# aurora/rootfs/etc/init.d/aurora
#!/sbin/openrc-run

name="aurora"
command=/sbin/agetty
command_args="--autologin root --noclear 38400 tty1 linux"
pidfile="/run/agetty.autologin.pid"
command_background=true
description="automatically login as root"

depend() {
    need localmount
    after bootmisc
}